# pkg/scripts/prerm.sh
#!/usr/bin/env bash

deb-systemd-helper stop haulage.service